# static/etc/services.d/vnc/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid bytesized /app/start_vnc
