# static/etc/services.d/xnovnc/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid bytesized /app/novnc/utils/launch.sh --vnc 127.0.0.1:5900
